# node-exporter.service
[Unit]
Description=Prometheus Node Exporter Service
After=network.target

[Service]
User=node_exporter
Group=node_exporter
Type=simple
ExecStart=/usr/local/bin/node_exporter

[Install]
<<<<<<< HEAD
WantedBy=multi-user.target
=======
WantedBy=multi-user.target
>>>>>>> 59822bb (node-exporter.service)
<<<<<<< HEAD
(local version)
=======
(remote version)
>>>>>>> origin/main


